# scripts/systemd/wifi-portal.service
[Unit]
Description=Wi-Fi Setup Captive Portal 
Wants=network.target 
After=network.target NetworkManager.service 

[Service]
Type=simple
ExecStart=/usr/local/sbin/wifi-portal.sh 
RemainAfterExit=yes
# The script runs until Wi-Fi connects successfully, then exits.

[Install]
WantedBy=multi-user.target
